# etc/cron.d/maas-gc
# Runs 'maas gc' every day at midnight
0 0 * * * root	maas gc &> /dev/null
